# etc/udev/rules.d/10-local.rules
# Set bluetooth power up
ACTION=="add", SUBSYSTEM=="bluetooth", RUN+="/usr/bin/bluetoothctl power on"
